# fs/usr/lib/onboard/hostapd-manager@ap0.path
[Unit]
Description=Trigger the hostapd-manager@ap0 unit when the network status changes

[Path]
PathChanged=/var/lib/onboard/network

[Install]
WantedBy=multi-user.target
